# templates/systemd.timer.j2
[Unit]
Description=Run minecraftoverview service

[Timer]
OnCalendar={{ minecraftoverviewer_on_calender }}
Persistent=true

[Install]
WantedBy=timers.target
